# sources/meta-compulab/recipes-bsp/wifi-start/files/wifi-start.service
[Unit]
Description=CompuLab WiFi Start Script

[Service]
Type=forking
ExecStart=/usr/sbin/wifi-start start

[Install]
WantedBy=multi-user.target
