# sources/meta-compulab/recipes-bsp/bt-start/files/bt-start.service
[Unit]
Description=CompuLab BLueThooth Start Script

[Service]
Type=forking
ExecStart=/usr/sbin/bt-start start

[Install]
WantedBy=multi-user.target
